# services.usersession/oneshot-user-late.service
[Unit]
Description=Oneshot stuff for user (late run)
After=user-session.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/oneshot --late

[Install]
WantedBy=post-user-session.target

