# services.usersession/oneshot-user.service
[Unit]
Description=Oneshot stuff for user
Requires=dbus.socket

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/oneshot

[Install]
WantedBy=pre-user-session.target

